# aggre-api.service
# File path: /etc/systemd/system/aggre-api.service
# sudo systemctl daemon-reload
# sudo systemctl start aggre-api
# sudo systemctl enable aggre-api
# sudo systemctl status aggre-api
[Unit]
Description=Aggre API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/aggre-api  # 注意修改路径
ExecStart=/path/to/aggre-api/aggre-api --port 3000 --log-dir /path/to/aggre-api/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
